# install/weathermap-updater.conf
description     "Weathermap SNMP device updater"

start on (started mysql)
stop on (stopped mysql)

respawn
respawn limit 2 5

chdir /opt/weathermap/scripts
exec ./update_devices.pl

